# highbase/rc-script
#!/bin/bash
#
# highbased
# this file is part of the highbase suite (named rc-script on the cvs tree)

# Copyright 2002 Fernando Ipar - fipar@acm.org / fipar@users.sourceforge.net

# This program is free software; you can redistribute it
# and/or modify it under the terms of the GNU General Public
# License as published by the Free Software Foundation;
# either version 2 of the License, or (at your option) any
# later version.
#
# This program is distributed in the hope that it will be
# useful, but WITHOUT ANY WARRANTY; without even the implied
# warranty of MERCHANTABILITY or FITNESS FOR A PARTICULAR
# PURPOSE.  See the GNU General Public License for more
# details.
#
# You should have received a copy of the GNU General Public
# License along with this program; if not, write to the Free
# Software Foundation, Inc., 59 Temple Place, Suite 330,
# Boston, MA 02111-1307 USA

### BEGIN INIT INFO
# Provides:         highbase
# Required-Start:    
# Required-Stop:     
# Should-Start:      
# Should-Stop:       
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the highbase daemon 
# Description:  highbase is a high availability clustering system for use 
# with the mysql database server. it uses mysql's builtin replication 
# features to achieve data redundance, fake to migrate a unique cluster IP 
# from master to slave in the event of a failure, and is highly customizable 
# and in active development
     
#                    
### END INIT INFO


# chkconfig: 2345 98 88
# description: highbase is a high availability clustering system for use \
# with the mysql database server. it uses mysql's builtin replication \
# features to achieve data redundance, fake to migrate a unique cluster IP \
# from master to slave in the event of a failure, and is highly customizable \
# and in active development
# processname: highbased
# config: /etc/highbase.conf
# pidfile: /var/run/highbase.pid

if [ -f "/etc/highbase.conf" ]; then
	. /etc/highbase.conf
elif [ -f "/etc/sysconfig/highbase" ]; then
	. /etc/sysconfig/highbase
elif [ -f "./highbase.conf" ]; then
	. ./highbase.conf
else
	echo "highbase not configured."
	exit 0
fi

export HIGHBASE_HOME

. /etc/profile

progname="$(basename $0)"
pidf="/var/run/highbase.pid"
logf="/var/log/highbase.log"
startedname="$HIGHBASE_HOME/configurator.sh"

. $HIGHBASE_HOME/role.include
. $HIGHBASE_HOME/compat.sh
. $HIGHBASE_HOME/common.sh


start()
{
	[ -f $pidf ] && echo "$progname already running (found pid file)">&2 && return 1
	echo -n "starting $progname: ">&2
	USER=$(cat $HIGHBASE_HOME/ssh_user)
	chown $USER.$USER $logf
	touch $pidf
	chown $USER.$USER $pidf
	encrypted=0
	grep ENCRYPTED /home/$USER/.ssh/id_dsa >/dev/null && {
		encrypted=1
		su - $USER -c "cd $HIGHBASE_HOME; env HIGHBASE_HOME=$HIGHBASE_HOME ./configurator.sh -o start-agent"	
	}
	nohup su - $USER -c "cd $HIGHBASE_HOME; env HIGHBASE_HOME=$HIGHBASE_HOME ./configurator.sh -p $pidf -e $encrypted>>$logf 2>&1" & >/dev/null 2>&1
	[ -f $pidf ] && {
		[ -d /proc/$(cat $pidf) ] && echo "OK">&2 || {
			echo "ERROR ($(cat $pidf) does not appear to be running)">&2
			return 1
		}
	} || echo "ERROR ($pidf not found after nohup)">&2 && return 1
	return 0
}

stop()
{
	[ -f $pidf ] || {
		echo "$progname does not appear to be running (pid fil not found)">&2
		return 1
	}

AGENT_SOCK=/tmp/highbase-ssh-agent.sock

[ -n "$(fuser $AGENT_SOCK)" ] && {
	kill $(fuser $AGENT_SOCK 2>&1|awk -F: '{print $2}') >/dev/null 2>&1
	/bin/usleep 200
	kill -9 $(fuser $AGENT_SOCK 2>&1|awk -F: '{print $2}') >/dev/null 2>&1
}

	echo -n "stopping $progname: ">&2
	su - $USER -c "cd $HIGHBASE_HOME; env HIGHBASE_HOME=$HIGHBASE_HOME ./configurator.sh -o shutdown-master" >/dev/null 2>&1
	pid=$(cat $pidf)
	echo "\$3 == $pid {print \$2}" > /tmp/awkscr.$$
	ps -ef|awk -f /tmp/awkscr.$$ |xargs kill >/dev/null 2>&1
	/bin/usleep 30
	ps -ef|awk -f /tmp/awkscr.$$ |xargs kill -9 >/dev/null 2>&1
	rm -f /tmp/awkscr.$$
	kill $(cat $pidf)	
	sleep 1
	[ -d /proc/$(cat $pidf) ] && {
		kill -9 $(cat $pidf) >/dev/null 2>&1
		echo "OK (had to kill -9)">&2
	} || echo "OK">&2
	rm -f $pidf
	return 0
}


restart() 
{
  	stop
	start
}	


status()
{
	[ -f $pidf ] && {
		echo "$progname is running ($(cat $pidf))">&2
	} || {
		echo "$progname is not running (no pidfile)">&2
	}
}

case "$1" in
  start)
  	start
	;;
  stop)
  	stop
	;;
  restart)
  	restart
	;;
  status)
  	status
	;;
  *)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
esac

exit $?
